# sh/prometheus/node_expoter批量安装脚本/node_exporter.service
[Unit]
Description=node_exporter
After=network.target

[Service]
Restart=on-failure
ExecStart=/usr/local/node_exporter/node_exporter/node_exporter

[Install]
WantedBy=multi-user.target